# templates/ai.ollama.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>ai.ollama.server</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string><string>-lc</string>
    <string>{{OLLAMA_BIN}} serve</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>OLLAMA_HOST</key><string>127.0.0.1:{{OLLAMA_PORT}}</string>
    <key>OLLAMA_FLASH_ATTENTION</key><string>{{OLLAMA_FLASH_ATTENTION}}</string>
    <key>OLLAMA_KV_CACHE_TYPE</key><string>{{OLLAMA_KV_CACHE_TYPE}}</string>
    <key>OLLAMA_KEEP_ALIVE</key><string>{{OLLAMA_KEEP_ALIVE}}</string>
    <key>OLLAMA_MAX_LOADED_MODELS</key><string>{{OLLAMA_MAX_LOADED_MODELS}}</string>
    <key>OLLAMA_NUM_PARALLEL</key><string>{{OLLAMA_NUM_PARALLEL}}</string>
  </dict>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>StandardOutPath</key><string>{{LOG_OLLAMA}}</string>
  <key>StandardErrorPath</key><string>{{LOG_OLLAMA_ERR}}</string>
</dict>
</plist>
