# templates/io.podman.machine.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>io.podman.machine</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string><string>-lc</string>
    <string>{{PODMAN_BIN}} machine start 2>/dev/null; exit 0</string>
  </array>
  <key>RunAtLoad</key><true/>
  <key>StandardOutPath</key><string>{{LOG_PODMAN}}</string>
  <key>StandardErrorPath</key><string>{{LOG_PODMAN}}</string>
</dict>
</plist>
